# config/systemd/user/mpd.service.d/override.conf
[Service]
RuntimeDirectory=mpd
RuntimeDirectoryMode=0755
