# config/systemd/user/ayugram.service
[Unit]
Description=AyuGram
After=graphical-session.target snixembed.service
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
ExecStart=/usr/sbin/AyuGram -startintray

[Install]
WantedBy=graphical-session.target
