# config/systemd/user/hyprsunset-reset.service
[Unit]
Description=Reset Hyprsunset to Default Profile

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyprsunset-reset-default.sh
