# config/systemd/user/loggy-user.path
[Unit]
Description=Watch for mouse flag file

[Path]
PathExists=/tmp/mouse_connected.flag

[Install]
WantedBy=graphical-session.target
