# config/systemd/user/profile-cleaner.timer
[Unit]
Description=Run profile-cleaner monthly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
