# config/systemd/user/profile-cleaner.service
[Unit]
Description=Run profile-cleaner

[Service]
Type=oneshot
ExecStart=/usr/sbin/profile-cleaner v
