# config/systemd/user/task-reminder.service
[Unit]
Description=Taskwarrior notify-send reminder

[Service]
Type=oneshot
ExecStart=/home/nick/Polka/local/bin/taskwarrior/task_notify.py
