# config/systemd/user/taskwarrior-notify.timer
[Unit]
Description=Run Task Reminder every 2 hours

[Timer]
OnBootSec=5m
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
