# config/systemd/user/taskwarrior-notify.service
[Unit]
Description=Taskwarrior notify-send reminder

[Service]
Type=oneshot
ExecStart=%h/.local/bin/taskwarrior/tw_notify.py
