# config/systemd/user/task-schedule.service
[Unit]
Description=Run Taskwarrior daily scheduling script

[Service]
Type=oneshot
ExecStart=%h/.local/bin/taskwarrior/task-schedule.py
