# config/systemd/user/wall.service
[Unit]
Description=Run Wall Script

[Service]
Type=oneshot
ExecStart=/home/nick/Polka/local/bin/wall/wall.py
